# cifs/rc.local
#!/bin/sh

sudo -u plex /etc/init.d/startup.sh

